# files/supervisord.conf
[supervisord]
nodaemon=true

[program:readthedocs]
directory=/www/readthedocs.org
command=python manage.py runserver 0.0.0.0:8000
autostart=true
autorestart=true